# src/crash-manager/crash-init.service.in
[Unit]
Description=Crash init
DefaultDependencies=no
After=local-fs.target

[Service]
Type=normal
ExecStart=@TZ_SYS_BIN@/@CRASH_INIT@.sh

[Install]
WantedBy=sysinit.target
